# devices/workstation/dot-config/systemd/user/sway-idle-inhibitor.service
[Unit]
Description=Inhibit sway idle while media is playing
PartOf=sway-session.target
After=sway-session.target

[Service]
ExecStart=%h/.local/bin/inhibit-on-playback.sh
Restart=on-failure

[Install]
WantedBy=sway-session.target

