# devices/workstation/dot-config/systemd/user/sway-session.target
[Unit]
Description=Sway-specific session setup
Wants=graphical-session.target
After=graphical-session.target
