# config/sway/dot-config/systemd/user/sway-wallpaper-update.service
[Unit]
Description=Update wallpaper and restart swaybg

[Service]
Type=oneshot
ExecStart=/bin/systemctl --user restart sway-wallpaper.service

