# config/sway/dot-config/systemd/user/swaybg.service
[Unit]
Description=swaybg wallpaper
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStartPre=/usr/bin/sh -lc 'for i in $(seq 1 100); do for s in "$XDG_RUNTIME_DIR"/wayland-*; do [ -S "$s" ] && exit 0; done; sleep 0.1; done; exit 1'
ExecStart=/usr/bin/swaybg --color 000000 --image %h/.config/sway/1.wallpaper --mode fill
Restart=on-failure
RestartSec=1s

[Install]
WantedBy=default.target

