# config/sway/dot-config/systemd/user/sway-wallpaper.timer
[Unit]
Description=Run wallpaper changer every 10 minutes

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
